# ansible/roles/install-startup-scripts/templates/zk-health-checker.service.j2
[Unit]
Description=Zookeeper Health checker
After=dcos-exhibitor.service

[Service]
Type=simple
Environment="PYTHONUNBUFFERED=1"
ExecStart={{ ansible_python_interpreter }} /opt/netsil/latest/bin/zk-health-checker.py

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
